# src/veil/environment/tunnel/tunnel-supervisord.cfg
[supervisord]
nodaemon = true
logfile=/tmp/tunnel-supervisord.log
pidfile=/tmp/tunnel-supervisord.pid
[inet_http_server]
port=*:9092
[program:socks_proxy]
command=veil environment tunnel socks-proxy-up
redirect_stderr=True
[program:redsocks]
command=veil environment tunnel redsocks-up
redirect_stderr=True